# rc.local
#!/bin/bash

/usr/bin/fanctrld &
exit 0
